# deploy/systemd/report-ai-portal.service.example
[Unit]
Description=RePORT AI Portal Streamlit UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=report-ai
Group=report-ai
WorkingDirectory=/opt/report-ai-portal
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=LOG_FORMAT=json
Environment=LOG_DIR=/var/log/report-ai-portal
Environment=STREAMLIT_SERVER_ADDRESS=127.0.0.1
Environment=STREAMLIT_SERVER_ENABLE_CORS=true
Environment=STREAMLIT_SERVER_ENABLE_XSRF_PROTECTION=true
Environment=REPORT_AI_PRODUCTION=1
Environment=REPORT_AI_REQUIRE_PHI_LOG_REDACTOR=1
Environment=REPORT_AI_AUTH_MODE=proxy
EnvironmentFile=/etc/default/report-ai-portal
ExecStart=/opt/report-ai-portal/.venv/bin/python main.py --web
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Linux hardening. Keep writable paths narrow and adjust them to the deployed
# data location before enabling this unit.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/report-ai-portal/data /opt/report-ai-portal/output /opt/report-ai-portal/tmp /var/log/report-ai-portal /var/lib/report-ai/.config/report_ai_portal
UMask=0077

[Install]
WantedBy=multi-user.target
